# deploy/systemd/vk-tutors-monitor.timer
[Unit]
Description=Monitor VK Tutors Bot

[Timer]
OnBootSec=2m
OnUnitActiveSec=2m
Unit=vk-tutors-monitor.service

[Install]
WantedBy=timers.target
